# deploy/deb/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop and disable the service
        if [ -d /run/systemd/system ]; then
            systemctl stop oxicrab.service || true
            systemctl disable oxicrab.service || true
            systemctl daemon-reload
        fi
        ;;
esac

#DEBHELPER#
